# docker/supervisor.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:slurm]
command=/usr/local/bin/slurmd -D
